# base/.config/systemd/user/weekly.service
[Unit]
Description=Runs weekly jobs

[Service]
Type=oneshot
ExecStart=%h/.local/share/scripts/maint
