# base/.local/share/services/photosync/run
#!/bin/sh

./bin/photosync-env photosync-watch
