# base/.config/autostart/redshift.desktop
[Desktop Entry]
Type=Application
Name=Redshift
Exec=/usr/bin/env redshift -O 4500
